# dot_config/systemd/user/clipse.service
[Unit]
Description=Clipboard manager

[Service]
Type=oneshot
ExecStart=/usr/bin/clipse -listen

[Install]
WantedBy=default.target
