# dot_config/systemd/user/openrgb-sakura.service
[Unit]
Description=Set OpenRGB Sakura theme

[Service]
Type=oneshot
ExecStart=/usr/bin/openrgb -p Sakura

[Install]
WantedBy=default.target
