# dot_config/systemd/user/udiskie.service
[Unit]
Description=Automount storage devices

[Service]
Type=oneshot
ExecStart=/usr/bin/udiskie

[Install]
WantedBy=default.target
